# ops/systemd/dealix-hermes.timer
[Unit]
Description=Run Dealix Hermes Agent OS every hour

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
Unit=dealix-hermes.service

[Install]
WantedBy=timers.target
